# debian/prerm
#!/bin/bash
set -e

# Script prerm para kioskapp
# Se ejecuta antes de desinstalar el paquete

SERVICE_NAME="printer-agent.service"

echo "Deteniendo servicios de KioskApp..."

# Detener y deshabilitar el servicio printer-agent
if systemctl is-active --quiet $SERVICE_NAME; then
    systemctl stop $SERVICE_NAME
fi

if systemctl is-enabled --quiet $SERVICE_NAME; then
    systemctl disable $SERVICE_NAME
fi

# Recargar systemd
systemctl daemon-reload

echo "Servicios detenidos"

exit 0





